# infra/data-services/docker/indexer.service
[Unit]
Description=Lucendex XRPL Indexer
After=docker.service lucendex-data-services.service
Requires=docker.service lucendex-data-services.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lucendex
EnvironmentFile=/opt/lucendex/.env
ExecStart=/opt/lucendex/indexer
Restart=always
RestartSec=10s
StandardOutput=append:/opt/lucendex/logs/indexer.log
StandardError=append:/opt/lucendex/logs/indexer.error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/lucendex/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
